# bluez-alsa/debian/postinst
#!/bin/bash

systemctl enable bluealsa bluealsa-aplay 
